# deploy/systemd/alpha-insight-watch-scheduler.service
[Unit]
Description=Alpha-Insight Watch Scheduler Worker
After=network-online.target alpha-insight-webhook-gateway.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/kkk/Project/Alpha-Insight
EnvironmentFile=/etc/alpha-insight/alpha-insight.env
Environment=PYTHONPATH=/home/kkk/Project/Alpha-Insight
ExecStart=/home/kkk/Project/Alpha-Insight/.venv/bin/python scripts/telegram_watch_scheduler.py --poll-interval-seconds 1 --batch-size 20
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/home/kkk/Project/Alpha-Insight/storage

[Install]
WantedBy=multi-user.target
